# courier.service.example
[Unit]
Description=courier Matrix relay bot
After=network-online.target
Wants=network-online.target

[Service]
User=courier
WorkingDirectory=/opt/courier
EnvironmentFile=/etc/courier.env
ExecStart=/opt/courier/courier
# The watchdog exits with code 70 on a stalled sync loop; always restart.
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
